# dockerfiles/build/www/configs/supervisor/supervisord.conf
[unix_http_server]
file=/run/supervisord.sock
; The username and password below are used only inside this file, they do not correspond to any OS user. But not
; providing them triggers a critical warning when starting supervisord.
; See: https://github.com/Supervisor/supervisor/issues/717
username = supervisord
password = TsuUeBYzatVm9vPMNJw42#&f#

[supervisord]
logfile=/dev/stdout
logfile_maxbytes=0
user=root

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///run/supervisord.sock
username = supervisord
password = TsuUeBYzatVm9vPMNJw42#&f#

[program:php-fpm]
command=php-fpm -F
autorestart=true
priority=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:nginx]
command=nginx -g "daemon off;"
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
