# dockerfiles/build/www/configs/supervisor/supervisord.conf
[unix_http_server]
file=/var/run/supervisord/supervisord.sock
chmod=0770
chown=www-data:www-data
; The username and password below are used only inside this file, they do not
; correspond to any OS user. But not providing them triggers a critical
; warning when starting supervisord.
; See: https://github.com/Supervisor/supervisor/issues/717
username = supervisord
password = TsuUeBYzatVm9vPMNJw42#&f#

[supervisord]
user=www-data
logfile=/var/log/supervisord/supervisord.log
pidfile=/var/run/supervisord/supervisord.pid

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///run/supervisord.sock
username = supervisord
password = TsuUeBYzatVm9vPMNJw42#&f#

[program:php-fpm]
command=php-fpm -F
autorestart=true
priority=5
stdout_logfile=/var/log/php/php-fpm-supervisord.log
redirect_stderr=true

[program:nginx]
command=nginx -g "daemon off;"
autorestart=true
priority=10
stdout_logfile=/var/log/nginx/nginx-supervisord.log
redirect_stderr=true
